# poweroff-at-9pm.service
[Unit]
Description=Scheduled Shutdown
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff
