# poweroff-at-9pm.service
[Unit]
Description=Scheduled Shutdown

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff
